# packaging/systemd/rust-junosmcp.service
[Unit]
Description=Rust Junos MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jmcp
Group=jmcp
ExecStart=/usr/local/bin/rust-junosmcp \
    --device-mapping /etc/jmcp/devices.json \
    --transport streamable-http \
    --host 127.0.0.1 \
    --port 30030 \
    --tokens-file /etc/jmcp/tokens.json \
    --device-lease-dir /var/lib/jmcp/device-leases \
    --inventory-readonly
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
UMask=0077
StateDirectory=jmcp
StateDirectoryMode=0750

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/jmcp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
